# service_startup/cloudflared-tunnel.service
[Unit]
Description=Cloudflare Tunnel for Dog Stream
After=network.target

[Service]
User=your-user
ExecStart=/usr/local/bin/cloudflared tunnel --config /home/your-user/.cloudflared/config.yml run dog-stream-tunnel
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target